# roles/pulldocker/templates/pulldocker.service.j2
[Unit]
Description=Pulldocker

[Service]
Type=oneshot
User=root
WorkingDirectory={{ pulldocker_path }}
ExecStart={{ pulldocker_path }}/bin/pulldocker --configuration pulldocker.yaml --quiet
